# extras/signman.socket
[Unit]
Description=Sign Manager Server socket
PartOf=signman.service

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=signman.service
